# spindownd.service
[Unit]
Description=Spindown Daemon

[Service]
SyslogIdentifier=spindownd

InaccessibleDirectories=/opt/ /boot/ /srv/ /var/
ProtectSystem=strict
ProtectHome=true
ProtectKernelModules=true
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=false

ExecStart=/usr/local/sbin/spindownd /dev/sdb /dev/sdc

Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
